# etc/cron.d/ledstriplight
# File: /etc/cron.d/ledstriplight
# Purpose: Automatically start/stop the `ledstriplight` systemd service at set times.
# Owner: root:root, Permissions: 644 (readable by root & cron daemon)

# File: /etc/cron.d/ledstriplight
# Purpose: Automatically start/stop the `ledstriplight` systemd service at set times.
# Owner: root:root, Permissions: 644 (readable by root & cron daemon)

# START
45 06 * * 1-5  root  /bin/systemctl start ledstriplight >> /tmp/log/ledstriplight-cron.log 2>&1   # weekdays (Mon–Fri)
30 07 * * 6,7  root  /bin/systemctl start ledstriplight >> /tmp/log/ledstriplight-cron.log 2>&1   # weekends (Sat & Sun)
05 17 * * *    root  /bin/systemctl start ledstriplight >> /tmp/log/ledstriplight-cron.log 2>&1   # daily

# STOP
30 09 * * 1-5  root  /bin/systemctl stop  ledstriplight >> /tmp/log/ledstriplight-cron.log 2>&1   # weekdays (Mon–Fri)
30 10 * * 6,7  root  /bin/systemctl stop  ledstriplight >> /tmp/log/ledstriplight-cron.log 2>&1   # weekends (Sat & Sun)
30 22 * * *    root  /bin/systemctl stop  ledstriplight >> /tmp/log/ledstriplight-cron.log 2>&1   # daily