# etc/systemd.d/ledstriplight.service
# LED Strip Light Controller Service
# This systemd service runs the LED Strip Light Controller Application as a daemon
# 
# Installation:
#   1. Copy this file to /etc/systemd/system/ledstriplight.service
#   2. Run: sudo systemctl daemon-reload
#   3. Run: sudo systemctl enable ledstriplight.service
#   4. Run: sudo systemctl start ledstriplight.service
#
# Management:
#   - Start:   sudo systemctl start ledstriplight.service
#   - Stop:    sudo systemctl stop ledstriplight.service
#   - Restart: sudo systemctl restart ledstriplight.service
#   - Status:  sudo systemctl status ledstriplight.service
#   - Logs:    sudo journalctl -u ledstriplight.service -f

[Unit]
Description=LED Strip Light Controller Service
After=network.target
# Ensure pigpio daemon is running before starting this service
Wants=pigpiod.service
After=pigpiod.service

[Service]
Type=simple
ExecStart=/usr/bin/python3 -u run.py profile
WorkingDirectory=/home/pi/raspberry/ledstriplight/src
User=pi
Group=pi
StandardOutput=append:/var/log/ledstriplight.log
StandardError=inherit
Restart=always
RestartSec=10
# Kill gracefully with SIGTERM, then SIGKILL after 30 seconds
KillMode=mixed
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
